# target/linux/imx25/base-files/etc/init.d/br_doze
#!/bin/sh /etc/rc.common
. /etc/br_doze.sh

START=55
echo "offmode=$offmode" 
echo "idletime=$idletime" 

start() {
  set_doze_shm	
}

